# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
user=root

; ─── Redis ───────────────────────────────────────────────────────────────────
[program:redis]
command=redis-server --appendonly yes --dir /home/data/redis --loglevel notice
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── DB migration (runs once at startup before web) ─────────────────────────
[program:migrate]
command=sh -c "sleep 5 && flask db upgrade"
directory=/app
environment=PYTHONPATH="/app",FLASK_APP="app.py"
autostart=true
autorestart=false
startsecs=0
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Admin init (runs once after migrate) ────────────────────────────────────
[program:admin-init]
command=sh -c "sleep 15 && python promote_admin_sqlalchemy.py"
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=false
startsecs=0
startretries=3
priority=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Seed default data (languages, catalogs, categories) ─────────────────────
[program:seed-data]
command=sh -c "sleep 20 && python add_llm_languages_table.py && python add_catalogs_table.py && python add_categories_table.py"
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=false
startsecs=0
startretries=3
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Celery worker ───────────────────────────────────────────────────────────
[program:celery]
command=celery -A celery_app.celery worker
    --loglevel=info
    --pool=solo
    -B
    --without-heartbeat
    --without-gossip
    --without-mingle
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
priority=20
; Wait for Redis and migrations to progress
startretries=5
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Gunicorn web server (Flask) ─────────────────────────────────────────────────────
[program:web]
command=gunicorn app:app
    --bind 0.0.0.0:%(ENV_PORT)s
    --workers %(ENV_WORKERS)s
    --timeout %(ENV_GUNICORN_TIMEOUT)s
    --access-logfile -
    --error-logfile -
    --log-level info
directory=/app
environment=PYTHONPATH="/app"
autostart=true
autorestart=true
priority=30
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── FastAPI server (v1.1 Admin Dashboard API) ─────────────────────────────────
[program:fastapi]
command=uvicorn main_fastapi:app
    --host 0.0.0.0
    --port 8001
    --log-level info
directory=/app
environment=PYTHONPATH="/app",VECTOR_STORE_PROVIDER="sqlite-vec"
autostart=true
autorestart=true
priority=31
startretries=5
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Web UI for supervisord (optional, disable if not needed) ────────────────
[inet_http_server]
port=127.0.0.1:9001
